# docker/puppet-master/puppet.service
[Unit]
Description=Puppet master

[Service]
ExecStart=/opt/puppet-master.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target